# racosync.service
[Unit]
Description=Racó sync attached files

[Service]
Type=oneshot
ExecStart=/usr/bin/python /home/rodrigo/dev/raco-api/racosync.py
